# packaging/openpropanel.service
[Unit]
Description=Open ProPanel — open-source server control panel
Documentation=https://github.com/jaskaranvirk1/openpropanel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/openpropanel --config /etc/openpropanel/config.json
Restart=on-failure
RestartSec=5

# Open ProPanel manages Apache, PHP-FPM and certbot and writes their config, so it
# runs privileged. StateDirectory gives it /var/lib/openpropanel with 0700 perms.
StateDirectory=openpropanel
StateDirectoryMode=0700

# Hardening that still lets the panel do its job (write /etc and /var config,
# run certbot/httpd/php-fpm, and provision system users). ProtectSystem=true
# keeps /usr and /boot read-only while leaving /etc + /var writable. We do NOT
# set ProtectHome — the panel provisions tenant home directories and may serve
# doc roots located under /home.
ProtectSystem=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
