# packaging/openpropanel.service
[Unit]
Description=Open ProPanel — open-source server control panel
Documentation=https://github.com/jaskaranvirk1/openpropanel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/openpropanel --config /etc/openpropanel/config.json
Restart=on-failure
RestartSec=5

# Open ProPanel manages Apache, PHP-FPM and certbot and writes their config, so it
# runs privileged. StateDirectory gives it /var/lib/openpropanel with 0700 perms.
StateDirectory=openpropanel
StateDirectoryMode=0700

# Hardening that still lets the panel do its job (write /etc and /var config,
# run certbot/httpd/php-fpm, and provision system users). ProtectSystem=true
# keeps /usr and /boot read-only while leaving /etc + /var writable. We do NOT
# set ProtectHome — the panel provisions tenant home directories and may serve
# doc roots located under /home.
ProtectSystem=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

# Drop the dangerous capabilities the panel never needs (inverted set, so the
# ones useradd/certbot/systemctl/httpd -t DO need — CHOWN, DAC_OVERRIDE, FOWNER,
# SETUID/SETGID, AUDIT_WRITE, … — stay available). A tighter positive allowlist
# is a follow-up once validated on a live AlmaLinux host.
CapabilityBoundingSet=~CAP_SYS_ADMIN CAP_SYS_MODULE CAP_SYS_RAWIO CAP_SYS_PTRACE CAP_SYS_BOOT CAP_SYS_TIME CAP_SYS_TTY_CONFIG CAP_NET_ADMIN CAP_MKNOD CAP_WAKE_ALARM CAP_AUDIT_CONTROL CAP_MAC_ADMIN CAP_MAC_OVERRIDE
# Hide other processes, block new namespaces, restrict socket families to the
# ones actually used, and refuse non-native syscall ABIs. (ProcSubset=pid is
# deliberately NOT set — it would hide /proc/meminfo,stat,loadavg and break the
# dashboard's resource stats.)
ProtectProc=invisible
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
